# hassio-wyze-gstreamer/root/etc/s6-overlay/s6-rc.d/cam2-snapshot/run
#!/usr/bin/with-contenv bashio
set -eu

bashio::log.info "Starting Cam2 Snapshot Stream"
# Create main config

RTSP_SIMPLE_SERVER_IP=$(bashio::config 'rtsp_simple_server_ip')
RTSP_SIMPLE_SERVER_PORT=$(bashio::config 'rtsp_simple_server_port')
CAM2_RTSP_SS_HI_RES_PATH=$(bashio::config 'cam2_rtsp_ss_hi_res_path')
CAM2_SNAPSHOT_NAME=$(bashio::config 'cam2_snapshot_name')

rm -f gstream-cam2-snapshot.log

tail --pid $$ -F gstream-cam2-snapshot.log &

gst-launch-1.0 rtspsrc protocols=tcp \
location=rtsp://${RTSP_SIMPLE_SERVER_IP}:${RTSP_SIMPLE_SERVER_PORT}/${CAM2_RTSP_SS_HI_RES_PATH} \
latency=0 ! rtph264depay ! avdec_h264 ! videorate ! video/x-raw,framerate=150/1001 ! \
jpegenc ! multifilesink location=/config/www/${CAM2_SNAPSHOT_NAME}.jpg -e > gstream-cam2-snapshot.log