# hassio-wyze-gstreamer/root/etc/s6-overlay/s6-rc.d/cam2-hi/run
#!/usr/bin/with-contenv bashio
set -eu

bashio::log.info "Starting Cam2 Hi Res Stream"
# Create main config

RTSP_SIMPLE_SERVER_IP=$(bashio::config 'rtsp_simple_server_ip')
RTSP_SIMPLE_SERVER_PORT=$(bashio::config 'rtsp_simple_server_port')
CAM2_IP=$(bashio::config 'cam2_ip')
CAM2_RTSP_PORT=$(bashio::config 'cam2_rtsp_port')
CAM2_USERNAME=$(bashio::config 'cam2_username')
CAM2_PASSWORD=$(bashio::config 'cam2_password')
CAM2_RTSP_SS_HI_RES_PATH=$(bashio::config 'cam2_rtsp_ss_hi_res_path')

rm -f gstream-cam2-hi.log

tail --pid $$ -F gstream-cam2-hi.log &

gst-launch-1.0 rtspsrc protocols=tcp \
location=rtsp://${CAM2_USERNAME}:${CAM2_PASSWORD}@${CAM2_IP}:${CAM2_RTSP_PORT}/video1_unicast \
latency=120 buffer-mode=3 connection-speed=2000 add-reference-timestamp-meta=true \
name=t t. ! queue ! \
capsfilter caps="application/x-rtp,media=video,height=1920,width=1080" ! \
queue ! rtph264depay ! h264parse ! h264timestamper ! queue ! \
rtspclientsink protocols=tcp latency=0 \
location=rtsp://${RTSP_SIMPLE_SERVER_IP}:${RTSP_SIMPLE_SERVER_PORT}/${CAM2_RTSP_SS_HI_RES_PATH} \
name=pay t. ! queue ! \
capsfilter caps="application/x-rtp,media=audio,clock-rate=16000,encoding-name=L16" ! \
rtpL16depay ! queue ! audioconvert ! queue ! \
volume volume=1.5 ! audioresample ! \
opusenc audio-type=2051 bandwidth=-1000 \
bitrate=64000 frame-size=40 ! queue ! \
pay. -e > gstream-cam2-hi.log
