# bumper/root/etc/s6-overlay/s6-rc.d/bumper/run
#!/usr/bin/with-contenv bashio
set -e
bashio::log.info "Starting Bumper"
# Create main config
BUMPER_ANNOUNCE_IP=$(bashio::config 'announce_ip')
TZ=$(bashio::config 'time_zone')
BUMPER_LISTEN=$(bashio::config 'announce_ip')
BUMPER_DEBUG=$(bashio::config 'debug_logs')
ENABLE_XMPP=True
ENABLE_MQTT=True

DISABLE_MQTT=$(bashio::config 'disable_mqtt_server')
DISABLE_XMPP=$(bashio::config 'disable_xmpp_server')

argline="--listen ${BUMPER_LISTEN} --announce ${BUMPER_ANNOUNCE_IP}"

if [ "${BUMPER_DEBUG}" = True ]
then
    argline+=" --debug"
fi

if [ "${DISABLE_XMPP}" = True ]
then
    argline+=" --disable-xmpp"
    ENABLE_XMPP=False
fi

if [ "${DISABLE_MQTT}" = True ]
then
    argline+=" --disable-mqtt"
    ENABLE_MQTT=False
fi

bashio::log.info "Announce IP is ${BUMPER_ANNOUNCE_IP}"
bashio::log.info "Listen IP is ${BUMPER_LISTEN}"
bashio::log.info "Debug is ${BUMPER_DEBUG}"
bashio::log.info "XMPP Server active is ${ENABLE_XMPP}"
bashio::log.info "MQTT Server active is ${ENABLE_MQTT}"
bashio::log.info "python command ran is:"
bashio::log.info "python3 -m bumper ${argline}"

cd /

exec python3 -m bumper ${argline}